# build/supervisor/builtin-webserver.conf
[program:builtin-webserver]
command=php -S 0.0.0.0:8080 -t /code/public
autostart=true
autorestart=true
priority=10
redirect_stderr=true
