# build/supervisor/php-fpm.conf
[program:php-fpm]
command=php-fpm --force-stderr --nodaemonize --fpm-config=/usr/local/etc/php-fpm.d/www.conf
autostart=true
autorestart=true
priority=5
redirect_stderr=true
user=USER-NAME
group=GROUP-NAME
